# conf/odst.conf
[program:blog]
command=pipenv run gunicorn -w 4 wsgi:app
directory=/root/odst
user=root
autostrat=true
autorestart=true
stopasgroup=true
killasgroup=true
